# packages/luci-app-360t7-hwaccel/data/etc/uci-defaults/95-360t7-hwaccel
#!/bin/sh

uci -q set firewall.@defaults[0].flow_offloading='1'
uci -q set firewall.@defaults[0].flow_offloading_hw='1'
uci -q commit firewall

exit 0
